# config/includes.chroot/etc/cron.hourly/ntpdate
#!/bin/sh

sudo /usr/sbin/ntpdate 0.ru.pool.ntp.org
